# Vbox-pet-lab/app/petclinic.service
[Unit]
Description=Spring Pet Clinic Application
After=network.target

[Service]
Type=simple
ExecStart=/app/start_app.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target